# etc/init.d/cron-openrc
#!/sbin/openrc-run

description="Starts the cron periodic command scheduler"
command="/usr/sbin/cron"
command_args="-f"
supervisor=supervise-daemon

depend() {
	need localmount
	use net
	provide cron
}
